# service/buildagent/run
#!/bin/bash
set -e
# echo "Setting serverUrl to '$TEAMCITY_SERVER'"
# sed -i -e "s/serverUrl=.*$/serverUrl=$(echo $TEAMCITY_SERVER | sed -e 's/[\/&]/\\&/g')/" /opt/buildAgent/conf/buildAgent.properties
#if [[ $TEAMCITY_OWN_ADDRESS ]]; then
#    echo "Setting ownAddress to '$TEAMCITY_OWN_ADDRESS'"
#    sed -i -e "s/#ownAddress=.*$/ownAddress=$TEAMCITY_OWN_ADDRESS/" /opt/buildAgent/conf/buildAgent.properties
#fi
echo "Starting buildagent..."
chown -R teamcity:root /opt/buildAgent
chown -R teamcity:root /data
exec /sbin/setuser teamcity /opt/buildAgent/bin/agent.sh run