# service/buildagent/run
#!/bin/bash
set -e

if [ -z "$TEAMCITY_SERVER" ]; then
    echo "TEAMCITY_SERVER variable not set, launch with -e TEAMCITY_SERVER=http://mybuildserver:myPort"
    exit 1
fi

if [ ! -d "$AGENT_DIR" ]; then
    echo "$AGENT_DIR doesn't exist pulling build-agent from server $TEAMCITY_SERVER";
    wget $TEAMCITY_SERVER/update/buildAgent.zip
    unzip -nd $AGENT_DIR buildAgent.zip
    rm buildAgent.zip
    chmod +x $AGENT_DIR/bin/agent.sh
fi

echo "Starting build-agent in $AGENT_DIR"
exec /sbin/setuser teamcity-agent $AGENT_DIR/bin/agent.sh run